# deploy/openbrain-telegram.service
[Unit]
Description=OpenBrain Telegram Bot
Documentation=https://github.com/windingriverholdings/openbrain
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=%h/repos/openbrain

EnvironmentFile=%h/repos/openbrain/.env

ExecStart=%h/repos/openbrain/bin/openbrain-telegram

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openbrain-telegram

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/repos/openbrain

[Install]
WantedBy=default.target
