# deploy/openbrain-watchd.service
[Unit]
Description=OpenBrain sandbox file bridge daemon
After=network.target postgresql.service
PartOf=openbrain-web.service

[Service]
Type=simple
WorkingDirectory=%h/repos/my-openbrain

ExecStart=%h/repos/my-openbrain/bin/openbrain-watchd

EnvironmentFile=%h/repos/my-openbrain/.env

Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/repos/my-openbrain %h/.openclaw

[Install]
WantedBy=default.target
