# deploy/openbrain-slack.service
[Unit]
Description=OpenBrain Slack Bot
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=%h/repos/openbrain

EnvironmentFile=%h/repos/openbrain/.env

ExecStart=%h/repos/openbrain/bin/openbrain-slack

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openbrain-slack

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/repos/openbrain

[Install]
WantedBy=default.target
